# rpm/cron.hourly
#! /bin/sh

runuser -c '/usr/libexec/resalloc-merge-hook-logs' - resalloc
